# docker/sipmedia/etc/profile.d/freeswitch.sh
alias fs_debug_on='/usr/bin/fs_cli -x "fsctl loglevel 7" && /usr/bin/fs_cli -x "sofia global siptrace on"'
alias fs_debug_off='/usr/bin/fs_cli -x "fsctl loglevel 3" && /usr/bin/fs_cli -x "sofia global siptrace off"'
